# packaging/sync-manager.service
[Unit]
Description=sync manager service startup script

[Service]
User=system
Group=system
SmackProcessLabel=sync-service
Type=simple
ExecStart=/usr/bin/sync-service
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
